# tools/com.tokyari.userprofilesync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.tokyari.userprofilesync</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string><string>-lc</string>
    <string>export PATH="$(/bin/ls -d "$HOME"/.nvm/versions/node/*/bin 2>/dev/null | sort -V | tail -1):/usr/local/bin:/opt/homebrew/bin:$PATH"; set -a; source /Users/oscardodds/projects/10koma-shukatsu/tools/.env.phase_c; set +a; /usr/bin/python3 /Users/oscardodds/projects/10koma-shukatsu/tools/sync_user_profiles.py</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/oscardodds/projects/10koma-shukatsu</string>
  <!-- 毎朝7:05。Mac が睡眠中ならlaunchdが起床時に実行(StartCalendarIntervalの標準挙動)。 -->
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>5</integer></dict>
  <!-- キャッチアップ: 電源OFF/再起動/ログイン等でロードされた時にも実行(7:05に走れなかった日を取り戻す)。同期は全件洗い替え=冪等。 -->
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/oscardodds/projects/10koma-shukatsu/logs/userprofilesync.out.log</string>
  <key>StandardErrorPath</key><string>/Users/oscardodds/projects/10koma-shukatsu/logs/userprofilesync.err.log</string>
</dict></plist>
