# tools/com.tokyari.companynews.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.tokyari.companynews</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>export PATH="$(/bin/ls -d "$HOME"/.nvm/versions/node/*/bin 2&gt;/dev/null|sort -V|tail -1):/usr/local/bin:/opt/homebrew/bin:$PATH"; exec /Users/oscardodds/oscar-ai/tokyari-pipeline/.venv/bin/python /Users/oscardodds/oscar-ai/tokyari-pipeline/scripts/crawl_company_news.py</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/Users/oscardodds/projects/10koma-shukatsu/logs/companynews.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/oscardodds/projects/10koma-shukatsu/logs/companynews.out.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>30</integer>
		<key>Weekday</key>
		<integer>1</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/oscardodds/projects/10koma-shukatsu</string>
</dict>
</plist>
